# deployment/example.fore-scheduler.service
# /etc/systemd/system/fore-scheduler.service
[Unit]
Description=Fore-Poster Scheduler
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/var/www/fore-poster
EnvironmentFile=/var/www/fore-poster/.env
Environment="PATH=/var/www/fore-poster/venv/bin"
ExecStart=/var/www/fore-poster/venv/bin/python fore_scheduler.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target